# ubuntu/16.04/etc/rc.local
#!/bin/bash

for f in /etc/rc.local.d/*; do
	if [ -x "$f" ]; then
		. "$f"
	fi
done

exit 0

# madneon 2014-12-16
