# Ubuntu/20.04/files/etc/cron.d/bkp
0 3 * * *   root    /opt/bin/bkp; /opt/bin/bkp-mirror
