# Ubuntu/16.04/files/etc/cron.daily/os-deploy
#!/bin/bash

/opt/lib/os-deploy/bin/deploy
